# pi/goodmorning-kiosk.service
[Unit]
Description=Good Morning Dashboard — Chromium Kiosk
After=graphical.target goodmorning-web.service
Wants=goodmorning-web.service

[Service]
Type=simple
User=goodmorning
Group=goodmorning
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-1
Environment=XDG_RUNTIME_DIR=/run/user/1001

# Wait for the web server to be ready before launching the browser
ExecStartPre=/bin/bash -c 'for i in $(seq 1 30); do curl -sf http://localhost/api/weather/ && break || sleep 1; done'

ExecStart=/usr/bin/chromium-browser \
    --kiosk \
    --noerrdialogs \
    --disable-translate \
    --no-first-run \
    --disable-infobars \
    --disable-session-crashed-bubble \
    --disable-features=TranslateUI \
    --check-for-update-interval=31536000 \
    --autoplay-policy=no-user-gesture-required \
    --password-store=basic \
    http://localhost

Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
